# defaults/bots/crontab
# SolidRusT defaults
45 *    * * *   admin   /bin/sh -c /root/solidrust.net/defaults/bots/45_logrotate.sh